# automation/systemd/patternpilot-sync.service
[Unit]
Description=Patternpilot watchlist automation

[Service]
Type=oneshot
WorkingDirectory=/path/to/patternpilot
ExecStart=/usr/bin/npm run automation:dispatch
ExecStartPost=/usr/bin/npm run automation:alerts
